# .config/autostart/.tablet_secondary-click-enabled.desktop
[Desktop Entry]
Type=Application
Exec=gsettings set org.gnome.desktop.a11y.mouse secondary-click-enabled "true"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=unity-integrated-menus
